# debian/postrm
#!/bin/bash

set -e

if [ "$1" = "purge" ]; then
   update-rc.d yunohost-firewall remove >/dev/null
fi

rm -f /etc/yunohost/installed

#DEBHELPER#

exit 0
